# services/dashcam-monitor/run
#!/bin/sh
exec setuidgid pi /home/pi/teslacam/src/dashcam-monitor.js
