# services/dashcam-monitor/run
#!/bin/sh
exec /root/teslacam/src/dashcam-monitor.js 
