# services/dashcam-monitor/run
#!/bin/sh
exec ~/teslacam/src/dashcam-monitor.js
