# services/dropbox-upload/run
#!/bin/sh
exec /root/teslacam/src/dropbox-upload.js 
