# services/dropbox-upload/run
#!/bin/sh
exec ~/teslacam/src/dropbox-upload.js
