# services/remote/run
#!/bin/sh
exec npm  --prefix ~/teslacam/src/remote run dev
