# services/remote/run
#!/bin/sh
exec /root/teslacam/src/remote && npm run dev
