# services/remote/run
#!/bin/sh
exec npm  --prefix /root/teslacam/src/remote run dev
