# services/remote/run
#!/bin/sh
exec setuidgid pi npm  --prefix /home/pi/teslacam/src/remote run dev
