# services/rotate-video/run
#!/bin/sh
exec setuidgid pi /home/pi/teslacam/src/rotate-video.js
