# installation/rc.local


if [ -e /boot/setup.txt ]
then
    /bin/bash /home/pi/teslacam/installation/postinstall.sh
    /bin/rm /boot/setup.txt
fi

