# installation/rc.local


if [ -e /boot/setup.txt]
then
    sh /home/pi/teslacam/installation/postinstall.sh
    rm /boot/setup.txt
fi

